# server/systemd/transcriber-cleanup.service
[Unit]
Description=Transcriber periodic cleanup
After=transcriber.service transcriber-bot.service

[Service]
Type=oneshot
ExecStart=/opt/transcriber/cleanup.sh
StandardOutput=journal
StandardError=journal
