# server/systemd/transcriber.service
[Unit]
Description=Transcriber audio to docx pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/transcriber
ExecStart=/opt/transcriber/venv/bin/python /opt/transcriber/transcriber.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/root
Environment=LLM_BASE_URL=https://your-llm-endpoint.example
Environment=LLM_API_KEY=YOUR_LLM_API_KEY

[Install]
WantedBy=multi-user.target
